# scripts/launchd/com.bmmpos.offline-pos.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.bmmpos.offline-pos</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>REPLACE_WITH_REPO_PATH/scripts/run_offline_pos.sh</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>OFFLINE_POS_ENV_FILE</key>
    <string>/Users/REPLACE_WITH_USER/Library/Application Support/BMM-POS/offline/offline-pos.env</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/REPLACE_WITH_USER/Library/Logs/BMM-POS/offline-pos.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/REPLACE_WITH_USER/Library/Logs/BMM-POS/offline-pos.error.log</string>
</dict>
</plist>
